# conf/crontab
@reboot /root/moss/switch/moss_quagga/init-quagga.sh